# deploy/hetzner/dashboard.service
[Unit]
Description=OpenClaw Dashboard (Next.js)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=openclaw
Group=openclaw
WorkingDirectory=/opt/openclaw/dashboard

# Shared runtime environment (do not hardcode secrets in unit files)
EnvironmentFile=/etc/openclaw/.env
Environment=NODE_ENV=production

ExecStart=/usr/bin/node /opt/openclaw/dashboard/node_modules/next/dist/bin/next start -p 3001

Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
